# profile.d/ssh.sh
#
# SSH Opts
#
export SSH_OPTS="-o StrictHostKeyChecking=no -o ForwardAgent=yes -o UserKnownHostsFile=/dev/null"
