# profile.d/ssh.sh
#
# SSH Opts
#
export SSH_OPTS="-A -o StrictHostKeyChecking=no"
