# etc/profile.d/alias.sh
#!/usr/bin/env bash

alias l='ls -lha'
alias ll='ls -lh'
alias vi=vim
alias pm-suspend='systemctl suspend'
