# etc/profile.d/alias.sh
#!/usr/bin/env bash

alias l='ls -lha'
alias ll='ls -lh'
alias vi=vim
alias ssh='ssh -XA'
alias wget='wget -c'

alias gerp=grep
alias gti=git
alias be='bundle exec'

alias pm-suspend='systemctl suspend'
