# etc/profile.d/aws_region.sh
#!/usr/bin/env bash
#
#
#
[[ -z ${AWS_REGION-} ]] && export AWS_REGION="ap-southeast-2"

