# etc/profile.d/postgres.sh
#
# POSTGRESQL VARS
#
export POSTGRES_IGNORED_DBS=(
    postgres
    rdsadmin
    template0
    template1
)
